# rootfst_cpio2/lib/udev/rules.d/90-usb.rules
KERNEL=="sd*[0-9]",SUBSYSTEMS=="usb", SUBSYSTEM=="block", ,RUN+="/etc/udev/mount.sh"
